# build/nginxphp/supervisord.conf
[supervisord]
nodaemon = true
logfile = /var/log/supervisord.log
logfile_maxbytes = 10MB
pidfile = /var/run/supervisord.pid

[program:nginx]
command = nginx
user = root
autostart = true
autorestart = true

[program:php-fpm5.6]
command = php-fpm5.6
user = root
autostart = true
autorestart = true
